# systemd/mnt-host.automount
[Unit]
Description=Automatically mount to host device for debugging

[Automount]
Where=/mnt/host

[Install]
WantedBy=multi-user.target